# vdir/.config/systemd/user/vdirsync.service
[Unit]
Description=CardDAV and CalDAV synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
