# linux-util/.config/autostart/wayland-instance.desktop
[Desktop Entry]
Name=Wayland Instance
Exec=wayland-instance
Terminal=false
Type=Application
NoDisplay=true
X-GNOME-Autostart-enabled=true
